# deploy/ark-panel.service
[Unit]
Description=ark-panel self-hosted OpenClaw session panel
After=network.target

[Service]
Type=simple
# Copy this example to ~/.config/systemd/user/ and replace USER/repository paths.
WorkingDirectory=/home/USER/awa/ark-panel
EnvironmentFile=/home/USER/.config/ark-panel/panel.env
Environment=PATH=/home/USER/.nvm/versions/node/v22.22.0/bin:/usr/local/bin:/usr/bin
ExecStart=/usr/bin/env node /home/USER/awa/ark-panel/dist/src/server/main.js
Restart=on-failure
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=15
UMask=0077

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadOnlyPaths=/home/USER/.openclaw/agents/claude/sessions /home/USER/.openclaw/agents/main/sessions
ReadWritePaths=/home/USER/.local/share/ark-panel /home/USER/.openclaw/agents/panel-runtime-claude/sessions /home/USER/.openclaw/agents/panel-runtime-main/sessions /home/USER/.openclaw/agents/panel-memory-claude/sessions /home/USER/.openclaw/agents/panel-memory-main/sessions /home/USER/claude/memory /home/USER/clawd/memory
# Confirmation refreshes OpenClaw's derived per-agent memory indexes. Keep real
# sessions read-only; only their separate agent database directories need writes.
ReadWritePaths=/home/USER/.openclaw/agents/claude/agent /home/USER/.openclaw/agents/main/agent /home/USER/.openclaw/agents/panel-runtime-claude/agent /home/USER/.openclaw/agents/panel-runtime-main/agent /home/USER/.openclaw/agents/panel-memory-claude/agent /home/USER/.openclaw/agents/panel-memory-main/agent
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=default.target
